# 24contest3.ini
[program:24contest_uwsgi]
command=/usr/Python2.7/bin/uwsgi --ini-paste /var/www/24contest.yosida95.com/production.ini
process_name=%(program_name)s
numprocs=1
autostart=true
user=www-data
stopsignal=INT
redirect_stderr=true
stdout_logfile=/var/www/24contest.yosida95.com/log/%(program_name)s.log
environment=PYTHON_EGG_CACHE=/var/www/24contest.yosida95.com/.cache
